# altservervirtmanager/vm-auto-pause.service
[Unit]
Description=Auto pause VM when RAM heavy processes run
After=network.target

[Service]
ExecStart=/usr/local/bin/vm-auto-pause.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
